# systemd user units/.config/systemd/user/mbsync.path
# -*- mode: conf-unix; -*-
[Unit]
Description=Trigger syncing after changes are made to .notmuch

[Path]
PathChanged=%h/Mail/.notmuch

[Install]
WantedBy=default.target
