# systemd user units/.config/systemd/user/mpd.socket.d/dont-listen-on-network-port.conf
[Socket]
ListenStream=
ListenStream=%t/mpd/socket
